# config/vagrant/foreman/export_templates/upstart/process.conf.erb
start on starting <%= app %>-<%= name %>
stop on stopping <%= app %>-<%= name %>
respawn

exec su - <%= user %> -c 'cd <%= engine.root %>; source .env; <%= process.command %> >> <%= log %>/<%=name%>-<%=num%>.log 2>&1'
